# deploy/com.tradesignal.us.open.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tradesignal.us.open</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Volumes/workspace/workspace/tradesignal/scripts/run_tradesignal_us_open.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Volumes/workspace/workspace/tradesignal</string>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>21</integer>
      <key>Minute</key>
      <integer>35</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>22</integer>
      <key>Minute</key>
      <integer>35</integer>
    </dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/tradesignal-us-open.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/tradesignal-us-open.stderr.log</string>
</dict>
</plist>
