# deploy/com.tradesignal.hk.open.mean_reversion.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.tradesignal.hk.open.mean_reversion</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/ABSOLUTE/PATH/TO/tradesignal/scripts/run_tradesignal_hk_open_mean_reversion.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/tradesignal</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>38</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/tradesignal-hk-open-mean-reversion.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/tradesignal-hk-open-mean-reversion.stderr.log</string>
</dict>
</plist>
